# systemd/luks-open-nas.service
[Unit]
Description=Open LUKS encrypted NAS device
DefaultDependencies=no
After=local-fs-pre.target
Before=local-fs.target cryptsetup.target
Conflicts=umount.target
ConditionPathExists=/dev/sda
ConditionPathExists=/root/.nas-keyfile

[Service]
Type=oneshot
RemainAfterExit=yes
# Open LUKS device if not already open
ExecStartPre=/bin/sh -c 'if ! /usr/sbin/cryptsetup status secure_nas_crypt >/dev/null 2>&1; then exit 0; else exit 1; fi || true'
ExecStart=/usr/sbin/cryptsetup luksOpen /dev/sda secure_nas_crypt --key-file /root/.nas-keyfile
# Close LUKS device on shutdown
ExecStop=/usr/sbin/cryptsetup close secure_nas_crypt
TimeoutSec=0

[Install]
WantedBy=local-fs-pre.target
